# v7.6.18/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:singlestoredb-studio]
command=/usr/bin/singlestoredb-studio
autostart=true
autorestart=true
